# deployment/src/scripts/fleetmanager-app.service
[Service]
WorkingDirectory=/opt/eve
ExecStart=/opt/eve/fleetmanager.jar
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fleetmanager-app
User=eve
Group=eve
Environment=NODE_ENV=production
[Install]
WantedBy=multi-user.target
